# mkosi.images/03-ncatomic-base/mkosi.extra/usr/lib/systemd/system/systemd-sysext.service.d/firstboot.conf
[Unit]
ConditionPathExists=|!/etc/extensions

[Service]
ExecStartPre=/usr/bin/bash -c '[[ -f /etc/verity.d/mkosi.crt ]] || { \
 mkdir -p /etc/verity.d; \
 /usr/bin/openssl x509 -in /boot/EFI/mkosi.der -out /etc/verity.d/mkosi.crt -outform PEM; \
}'
ExecStartPre=/usr/bin/bash -c '[[ -d /etc/extensions ]] || { \
 mkdir -p /etc/extensions; \
 for ext in /usr/share/embedded/extensions/*.raw; \
 do \
  [[ "$${ext}" == "/usr/share/embedded/extensions/*.raw" ]] || ln -s "$${ext}" "/etc/extensions/$(basename "$${ext%_*}.raw")"; \
 done; \
}'
