# configs/ncp/nc_aio/root_overlay/usr/lib/systemd/system/ncp-activation.service
[Unit]
Description=NCP provisioning

[Service]
Type=oneshot
ConditionPathExists=!/etc/ncp/all-in-one/values.yaml
ExecStart=/opt/ncp/ncp-activation.sh