# mkosi.images/01-base/mkosi.extra/usr/lib/systemd/system/systemd-confext.service.d/mutable.conf
[Service]
ExecStart=
ExecStart=systemd-confext refresh --mutable=auto
ExecReload=
ExecReload=systemd-confext refresh --mutable=auto