# mkosi.images/71-ncatomic-sysext/mkosi.extra/usr/lib/systemd/system/credentials-setup.service
[Unit]

[Service]
Type=oneshot
ExecStart=/usr/bin/systemd-creds setup

[Install]
WantedBy=ncatomic.target