# mkosi.images/01-base/mkosi.extra/usr/lib/systemd/system/ssh.service.d/setup_service.conf
[Unit]
Requires=sshd_setup.service
After=sshd_setup.service
