# mkosi.images/71-ncatomic-sysext/mkosi.extra/usr/lib/systemd/system/extensions.service.d/ncatomic.conf
[Service]
ExecStartPre=systemctl enable ncatomic.target nca-web.service nca-occd.service nca-caddy.service nca-unlock.service nca-unlock-from-credentials.service nextcloud-all-in-one.service credentials-setup.service
ExecStartPre=systemctl start ncatomic.target